# template.management_app.service
[Unit]
Description=Management APP
After=syslog.target

[Service]
WorkingDirectory=%APP_DIR%
ExecStart=%PYTHON_PATH% -m app.main -mc main.ini -cc cloud.ini -nc nginx.conf
Restart=always
KillSignal=SIGQUIT
Type=simple
StandardError=file:/var/log/management_app.log
NotifyAccess=all

[Install]
WantedBy=multi-user.target
